# deploy/civics-congress-sync.service
[Unit]
Description=Sync recent federal bills from Congress.gov into Civics Radar
After=network-online.target civics-api.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/civics
EnvironmentFile=-/etc/civics.env
Environment=CIVICS_DB=/opt/civics/data/civics.db
ExecStart=/opt/civics/.venv/bin/python -m civics_app.sync_congress --limit 25
User=root
