# deploy/civics-db-backup.timer
[Unit]
Description=Back up Civics Radar daily

[Timer]
OnCalendar=daily
RandomizedDelaySec=30m
Persistent=true
Unit=civics-db-backup.service

[Install]
WantedBy=timers.target
