# deploy/civics-legiscan-sync.service
[Unit]
Description=Sync state bills from LegiScan into Civics Radar
After=network-online.target civics-api.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/civics
EnvironmentFile=-/etc/civics.env
Environment=CIVICS_DB=/opt/civics/data/civics.db
ExecStart=/opt/civics/.venv/bin/python -m civics_app.sync_legiscan --state MO --limit 25
User=civics
Group=civics
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
ReadWritePaths=/opt/civics/data
CapabilityBoundingSet=
AmbientCapabilities=
