# deploy/civics-notification-worker.service
[Unit]
Description=Deliver queued Civics Radar notifications
After=network-online.target civics-api.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/civics
EnvironmentFile=-/etc/civics.env
Environment=CIVICS_DB=/opt/civics/data/civics.db
ExecStart=/opt/civics/.venv/bin/python -m civics_app.sync_notifications --limit 100
User=civics
Group=civics
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
ReadWritePaths=/opt/civics/data
CapabilityBoundingSet=
AmbientCapabilities=
