# deploy/civics-db-backup.service
[Unit]
Description=Back up the Civics Radar SQLite database
After=civics-api.service

[Service]
Type=oneshot
WorkingDirectory=/opt/civics
EnvironmentFile=-/etc/civics.env
Environment=CIVICS_DB=/opt/civics/data/civics.db
Environment=CIVICS_BACKUP_DIR=/opt/civics/backups
ExecStart=/opt/civics/scripts/backup_civics_db.sh
User=civics
Group=civics
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
ReadWritePaths=/opt/civics/data /opt/civics/backups
CapabilityBoundingSet=
AmbientCapabilities=
